# upstart/tippr-consumers-start.conf
description "start up queue consumers"

task

start on runlevel [2345] or tippr-start

script
    . /etc/default/tippr
    manage-consumers
end script
