# upstart/tippr-boot.conf
description "start tippr on boot"

task
start on runlevel [2345]

exec initctl emit tippr-start
