# upstart/tippr-consumers-restart.conf
description "restart queue consumers"

task

start on tippr-restart or tippr-kill

script
    . /etc/default/tippr
    manage-consumers
end script
